# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null

[program:api]
command=uvicorn api.main:app --host 0.0.0.0 --port 8000
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true

[program:mcp]
command=python -m tether_mcp.server --sse
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true

[program:bot]
command=bash -c "if [ -z \"$TELEGRAM_WEBHOOK_URL\" ]; then until curl -sf http://localhost:8000/api/health > /dev/null; do sleep 2; done && exec python -m bot.message_handler; else echo '[bot] Webhook mode active — polling process not needed' && exec sleep infinity; fi"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=true
startsecs=10
