# systemd/tether-sync.service
[Unit]
Description=Tether Sync Worker
After=network.target

[Service]
Type=simple
User=toast
WorkingDirectory=/home/toast/tether
EnvironmentFile=-/home/toast/.tether-config/env
ExecStart=/home/toast/tether/.venv/bin/python -m sync
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
